# docker/app/supervisor-frontend-web/conf.d/keep-alive.conf
[program:keep_alive]
command=/bin/sh -c "while true; do echo 'run in dev mode' > /dev/null 2>&1; sleep 3600; done"
autostart=true
autorestart=true
stdout_logfile=/dev/null
stderr_logfile=/dev/null